# roles/ai_processing/templates/ai-worker.service.j2
[Unit]
Description=AI Redis Worker
After=network.target

[Service]
ExecStart=/usr/bin/python3 /opt/ai/redis_worker.py
Restart=always
User=ansible
WorkingDirectory=/opt/ai

[Install]
WantedBy=multi-user.target
